# Letux/etc/systemd/system/letux-panel-setup.service
[Unit]
Description=Setup Letux Panel config
Before=nodm.service

[Service]
Type=oneshot
ExecStart=/etc/X11/xinit/letux-panel-setup

[Install]
WantedBy=graphical.target
