# Letux/etc/systemd/system/tiler.service
[Unit]
Description=OMAP5 Tiler fix

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/root/tiler.sh

[Install]
Before=graphical.target
